# crontab
server:
*/5 * * * * /root/.iptell/listip >/dev/null 2>&1

client:
*/10 * * * * /home/will/.iptell/iptell > /dev/null 2?&1
