# crontab
server:
*/1 * * * * sleep 10; /server-home-here/.iptell/listip >/dev/null 2>&1

client:
*/1 * * * * /your-home-here/.iptell/iptell > /dev/null 2?&1
@reboot nohup /your-home-here/.iptell/usagelog &
